# packages/pi-base/pi_base-0.0.34.tar.gz/pi_base-0.0.34/pi_base/common/pkg/etc/systemd/system/hdmi-sound-rpi4-c1.target
[Unit]
Description=Sound via HDMI on RPI4 (card1)
RefuseManualStart=yes
Wants=hdmi-sound.target
Before=hdmi-sound.target

After=alsa-restore.service sound.target 
Requires=alsa-restore.service sound.target 

After=sys-devices-platform-soc-fef00700.hdmi-sound-card1-controlC1.device
Wants=sys-devices-platform-soc-fef00700.hdmi-sound-card1-controlC1.device
ConditionPathExists=|/sys/devices/platform/soc/fef00700.hdmi/sound/card1/controlC1/

[Install]
WantedBy=hdmi-sound.target
