# packages/psyserver/psyserver-0.4.3.tar.gz/psyserver-0.4.3/psyserver/example/psyserver.service
[Unit]
Description=PsyServer Service.

[Service]
# pyseserver init will place real paths here
ExecStart=/path/to/python /path/to/psyserver_package run /path/to/psyserver_dir
Environment=PYTHONUNBUFFERED=1
Restart=on-failure

[Install]
WantedBy=default.target

