# packages/pi-base/pi_base-0.0.35.tar.gz/pi_base-0.0.35/pi_base/common/pkg/etc/systemd/system/hdmi-sound.target
[Unit]
Description=Sound via HDMI
RefuseManualStart=yes


[Install]
Also=hdmi-sound-rpi3.target
Also=hdmi-sound-rpi4-c1.target
Also=hdmi-sound-rpi4-c2.target
